# daemon/claude-usage-daemon-usb.service
[Unit]
Description=Claude Usage Tracker Daemon (USB)

[Service]
Type=simple
ExecStart=DAEMON_PATH
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
